# utils/usrp-signal-manager.service
[Unit]
Description=USRP Signal Manager Service
After=network.target uhd-host.service
Wants=network-online.target

[Service]
User=your_username  # 替换为你的用户名
Group=your_group    # 替换为你的用户组
WorkingDirectory=/path/to/usrp-signal-tool  # 替换为项目根目录
ExecStart=/path/to/python3 app.py  # 替换为 Python 路径（虚拟环境路径或系统 Python）
Environment="PYTHONUNBUFFERED=1"
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target